# hosts/jenkins/etc/cron.d/dsa
# m h dom mon dow (0|7=sun,1=mon)    command

#
# cron-job for jenkins.debian.net
#

MAILTO=root

0 1 * * * nobody /usr/bin/chronic /usr/local/bin/dsa-check-running-kernel
2 1,13 * * * nobody /usr/bin/chronic /usr/local/bin/dsa-check-packages
0 0 * * * root mkdir -p /etc/apache2/ssl/debsso && /usr/local/bin/update-debsso-ca --destdir /etc/apache2/ssl/debsso
